# playbooks/roles/base_buri/files/etc/systemd/system/ephemeral-raid.service
Unit]
Description=Automagically formats ebs & ephemeral disks and creates defined symlinks
After=network.target

[Service]
ExecStart=/usr/local/bin/buri-drives.sh

[Install]
WantedBy=multi-user.target
