# playbooks/roles/cassandra/templates/cassandra-storage.conf
# cassandra-storage - Setup EBS storage for Cassandra usage.
#
# The first non-root EBS volume will be configured as the commit log disk.
# If there is only one EBS volume, data will be stored there, too.
# If there is a second EBS volume, data will be stored on it.
# If there are three or more EBS volumes, the 2nd and subsequent volumes
# will be configured in RAID 0 using mdadm.
#
# All volumes are configured with xfs.
# To get the most out of the GP2 EBS IOPS, a 256 KB block size is used.
# It's also better to use multiple volumes to take advantage of the 3000
# burstable IOPS per volume.

description "Setup storage for Cassandra usage"

start on stopped ephemeral-mkdirs
console log

script
    buri_configured_file=/var/local/buri_cassandra_storage

    if [ -e $buri_configured_file ]; then
      echo "Detected Cassandra storage already setup, exiting"
      exit 0
    fi

    METADATA_URL_BASE="http://169.254.169.254/2012-01-12"

    ebs_drives=""
    ebs_count=0
    ebses=$(curl --silent $METADATA_URL_BASE/meta-data/block-device-mapping/ | grep ebs)

    for e in $ebses; do
        echo "Probing $e ..."
        device_name=$(curl --silent $METADATA_URL_BASE/meta-data/block-device-mapping/$e)
        device_name=$(echo $device_name | sed "s/sd/xvd/")
        device_path="/dev/$device_name"
        if [ -b $device_path ]; then
            ebs_drives="$ebs_drives $device_path"
            ebs_count=$(($ebs_count+1))
        fi
    done

    case "${ebs_count}" in
        0)
            echo "No EBS volumes detected."
            ;;
        1)
            echo "Single EBS volume detected."
            mkfs.xfs -f -b size=4k -s size=4k -d su=256k,sw=1 ${device_path}
            mount -t xfs -o noatime ${device_path} /cassandra
            sed -i "#${device_path}#d" /etc/fstab
            echo "${device_path} /cassandra xfs noatime 0 0" | tee -a /etc/fstab
            mkdir /cassandra/data /cassandra/saved_caches /cassandra/commitlog
            chown {{ cassandra_user }}:{{ cassandra_group }} /cassandra/data /cassandra/saved_caches /cassandra/commitlog
            ;;
        2)
            echo "Two EBS volumes detected. Using first for commitlog, the second for data."

            mkdir -p /cassandra_data /cassandra_commitlog

            first=$(echo ${ebs_drives} | cut -d' ' -f1)
            second=$(echo ${ebs_drives} | cut -d' ' -f2-)

            mkfs.xfs -f -b size=4k -s size=4k -d su=256k,sw=1 ${second}
            mount -t xfs -o noatime ${second} /cassandra_data
            sed -i "#${second}#d" /etc/fstab
            echo "${second} /cassandra_data xfs noatime 0 0" | tee -a /etc/fstab

            mkfs.xfs -f -b size=4k -s size=4k -d su=256k,sw=1 ${first}
            mount -t xfs -o noatime ${first} /cassandra_commitlog
            sed -i "#${first}#d" /etc/fstab
            echo "${first} /cassandra_commitlog xfs noatime 0 0" | tee -a /etc/fstab

            # delete existing folders
            for i in /cassandra/data /cassandra/saved_caches /cassandra/commitlog ; do
                if [ -e $i ] ; then
                    rm -r $i
                fi
            done

            mkdir /cassandra_data/data /cassandra_data/saved_caches /cassandra_commitlog/commitlog
            chown {{ cassandra_user }}:{{ cassandra_group }} /cassandra_data/data /cassandra_data/saved_caches /cassandra_commitlog/commitlog
            ln -s /cassandra_commitlog/commitlog /cassandra/commitlog
            ln -s /cassandra_data/data /cassandra/data
            ln -s /cassandra_data/saved_caches /cassandra/saved_caches
            ;;
        *)
            echo "Several (3+) EBS volumes detected. Using first for commitlog, remainder for data."

            mkdir -p /cassandra_data /cassandra_commitlog

            first=$(echo ${ebs_drives} | cut -d' ' -f1)
            remaining=$(echo ${ebs_drives} | cut -d' ' -f2-)
            raid_count=$((${ebs_count} - 1))

            mdadm --create --verbose /dev/md1 --level=0 -c256 --raid-devices=${raid_count} ${remaining}
            echo DEVICE ${remaining} | tee /etc/mdadm/mdadm.conf
            mdadm --detail --scan | tee -a /etc/mdadm/mdadm.conf

            mkfs.xfs -f -b size=4k -s size=4k -d su=256k,sw=${raid_count} /dev/md1
            mount -t xfs -o noatime /dev/md1 /cassandra_data
            sed -i "#/dev/md1#d" /etc/fstab
            echo "/dev/md1 /cassandra_data xfs noatime 0 0" | tee -a /etc/fstab

            mkfs.xfs -f -b size=4k -s size=4k -d su=256k,sw=1 ${first}
            mount -t xfs -o noatime ${first} /cassandra_commitlog
            sed -i "#${first}#d" /etc/fstab
            echo "${first} /cassandra_commitlog xfs noatime 0 0" | tee -a /etc/fstab

            # delete existing folders
            for i in /cassandra/data /cassandra/saved_caches /cassandra/commitlog ; do
                if [ -e $i ] ; then
                    rm -r $i
                fi
            done

            mkdir /cassandra_data/data /cassandra_data/saved_caches /cassandra_commitlog/commitlog
            chown {{ cassandra_user }}:{{ cassandra_group }} /cassandra_data/data /cassandra_data/saved_caches /cassandra_commitlog/commitlog
            ln -s /cassandra_commitlog/commitlog /cassandra/commitlog
            ln -s /cassandra_data/data /cassandra/data
            ln -s /cassandra_data/saved_caches /cassandra/saved_caches

            update-initramfs -u
            ;;
    esac

    touch $buri_configured_file
    echo "Cassandra storage configured."
    exit 0
end script
