# playbooks/roles/nat/files/natmon.conf
description "Nat monitoring"
author "Marco Monaco"

start on (started networking)
stop on shutdown

respawn
respawn limit 99 5

script
    exec /usr/local/bin/natmon.sh
end script
